# scripts/k8s-node/charts/aipub/templates/aipub-ha.service.j2
[Unit]
Description=Aipub HA Service

[Service]
Type=simple
User=root
ExecStart={{ exec_start }}
Restart=on-failure

[Install]
WantedBy=default.target
